# tags/attack.t1053.002/index.xml
<?xml version="1.0" encoding="utf-8" standalone="yes"?>
<rss version="2.0" xmlns:atom="http://www.w3.org/2005/Atom">
  <channel>
    <title>attack.t1053.002 on My New Hugo Site</title>
    <link>https://bradleyjkemp.dev/sigmadoc/tags/attack.t1053.002/</link>
    <description>Recent content in attack.t1053.002 on My New Hugo Site</description>
    <generator>Hugo -- gohugo.io</generator>
    <language>en-us</language>
    <lastBuildDate>Mon, 10 Jan 2022 18:04:43 +0200</lastBuildDate><atom:link href="https://bradleyjkemp.dev/sigmadoc/tags/attack.t1053.002/index.xml" rel="self" type="application/rss+xml" />
    <item>
      <title>Remote Schedule Task Lateral Movement via ATSvc</title>
      <link>https://bradleyjkemp.dev/sigmadoc/rules/application/rpc_firewall/rpc_firewall_atsvc_lateral_movement.yml/</link>
      <pubDate>Mon, 10 Jan 2022 18:04:43 +0200</pubDate>
      
      <guid>https://bradleyjkemp.dev/sigmadoc/rules/application/rpc_firewall/rpc_firewall_atsvc_lateral_movement.yml/</guid>
      <description>&lt;p&gt;Detects remote RPC calls to create or execute a scheduled task via ATSvc&lt;/p&gt;</description>
    </item>
    
    <item>
      <title>Remote Schedule Task Lateral Movement via ITaskSchedulerService</title>
      <link>https://bradleyjkemp.dev/sigmadoc/rules/application/rpc_firewall/rpc_firewall_itaskschedulerservice_lateral_movement.yml/</link>
      <pubDate>Mon, 10 Jan 2022 18:04:43 +0200</pubDate>
      
      <guid>https://bradleyjkemp.dev/sigmadoc/rules/application/rpc_firewall/rpc_firewall_itaskschedulerservice_lateral_movement.yml/</guid>
      <description>&lt;p&gt;Detects remote RPC calls to create or execute a scheduled task&lt;/p&gt;</description>
    </item>
    
    <item>
      <title>Remote Schedule Task Lateral Movement via SASec</title>
      <link>https://bradleyjkemp.dev/sigmadoc/rules/application/rpc_firewall/rpc_firewall_sasec_lateral_movement.yml/</link>
      <pubDate>Mon, 10 Jan 2022 18:04:43 +0200</pubDate>
      
      <guid>https://bradleyjkemp.dev/sigmadoc/rules/application/rpc_firewall/rpc_firewall_sasec_lateral_movement.yml/</guid>
      <description>&lt;p&gt;Detects remote RPC calls to create or execute a scheduled task via SASec&lt;/p&gt;</description>
    </item>
    
    <item>
      <title>Interactive AT Job</title>
      <link>https://bradleyjkemp.dev/sigmadoc/rules/windows/process_creation/proc_creation_win_interactive_at.yml/</link>
      <pubDate>Mon, 28 Oct 2019 11:59:49 +0100</pubDate>
      
      <guid>https://bradleyjkemp.dev/sigmadoc/rules/windows/process_creation/proc_creation_win_interactive_at.yml/</guid>
      <description>&lt;p&gt;Detect an interactive AT job, which may be used as a form of privilege escalation.&lt;/p&gt;</description>
    </item>
    
    <item>
      <title>MITRE BZAR Indicators for Execution</title>
      <link>https://bradleyjkemp.dev/sigmadoc/rules/network/zeek/zeek_dce_rpc_mitre_bzar_execution.yml/</link>
      <pubDate>Mon, 01 Jan 0001 00:00:00 +0000</pubDate>
      
      <guid>https://bradleyjkemp.dev/sigmadoc/rules/network/zeek/zeek_dce_rpc_mitre_bzar_execution.yml/</guid>
      <description>&lt;p&gt;Windows DCE-RPC functions which indicate an execution techniques on the remote system. All credit for the Zeek mapping of the suspicious endpoint/operation field goes to MITRE&lt;/p&gt;</description>
    </item>
    
    <item>
      <title>Remote Task Creation via ATSVC Named Pipe</title>
      <link>https://bradleyjkemp.dev/sigmadoc/rules/windows/builtin/security/win_atsvc_task.yml/</link>
      <pubDate>Mon, 01 Jan 0001 00:00:00 +0000</pubDate>
      
      <guid>https://bradleyjkemp.dev/sigmadoc/rules/windows/builtin/security/win_atsvc_task.yml/</guid>
      <description>&lt;p&gt;Detects remote task creation via at.exe or API interacting with ATSVC namedpipe&lt;/p&gt;</description>
    </item>
    
    <item>
      <title>Remote Task Creation via ATSVC Named Pipe - Zeek</title>
      <link>https://bradleyjkemp.dev/sigmadoc/rules/network/zeek/zeek_smb_converted_win_atsvc_task.yml/</link>
      <pubDate>Mon, 01 Jan 0001 00:00:00 +0000</pubDate>
      
      <guid>https://bradleyjkemp.dev/sigmadoc/rules/network/zeek/zeek_smb_converted_win_atsvc_task.yml/</guid>
      <description>&lt;p&gt;Detects remote task creation via at.exe or API interacting with ATSVC namedpipe&lt;/p&gt;</description>
    </item>
    
  </channel>
</rss>
