# tags/attack.t1053/index.xml
<?xml version="1.0" encoding="utf-8" standalone="yes"?><rss version="2.0" xmlns:atom="http://www.w3.org/2005/Atom"><channel><title>attack.t1053 on My New Hugo Site</title><link>https://bradleyjkemp.dev/sigmadoc/tags/attack.t1053/</link><description>Recent content in attack.t1053 on My New Hugo Site</description><generator>Hugo -- gohugo.io</generator><language>en-us</language><lastBuildDate>Mon, 10 Jan 2022 18:04:43 +0200</lastBuildDate><atom:link href="https://bradleyjkemp.dev/sigmadoc/tags/attack.t1053/index.xml" rel="self" type="application/rss+xml"/><item><title>Remote Schedule Task Lateral Movement via ATSvc</title><link>https://bradleyjkemp.dev/sigmadoc/rules/application/rpc_firewall/rpc_firewall_atsvc_lateral_movement.yml/</link><pubDate>Mon, 10 Jan 2022 18:04:43 +0200</pubDate><guid>https://bradleyjkemp.dev/sigmadoc/rules/application/rpc_firewall/rpc_firewall_atsvc_lateral_movement.yml/</guid><description>&lt;p>Detects remote RPC calls to create or execute a scheduled task via ATSvc&lt;/p></description></item><item><title>Remote Schedule Task Lateral Movement via ITaskSchedulerService</title><link>https://bradleyjkemp.dev/sigmadoc/rules/application/rpc_firewall/rpc_firewall_itaskschedulerservice_lateral_movement.yml/</link><pubDate>Mon, 10 Jan 2022 18:04:43 +0200</pubDate><guid>https://bradleyjkemp.dev/sigmadoc/rules/application/rpc_firewall/rpc_firewall_itaskschedulerservice_lateral_movement.yml/</guid><description>&lt;p>Detects remote RPC calls to create or execute a scheduled task&lt;/p></description></item><item><title>Remote Schedule Task Lateral Movement via SASec</title><link>https://bradleyjkemp.dev/sigmadoc/rules/application/rpc_firewall/rpc_firewall_sasec_lateral_movement.yml/</link><pubDate>Mon, 10 Jan 2022 18:04:43 +0200</pubDate><guid>https://bradleyjkemp.dev/sigmadoc/rules/application/rpc_firewall/rpc_firewall_sasec_lateral_movement.yml/</guid><description>&lt;p>Detects remote RPC calls to create or execute a scheduled task via SASec&lt;/p></description></item><item><title>Suspicious Scheduled Task Write to System32 Tasks</title><link>https://bradleyjkemp.dev/sigmadoc/rules/windows/file_event/file_event_win_susp_task_write.yml/</link><pubDate>Tue, 16 Nov 2021 17:30:47 +0100</pubDate><guid>https://bradleyjkemp.dev/sigmadoc/rules/windows/file_event/file_event_win_susp_task_write.yml/</guid><description>&lt;p>Detects the creation of tasks from processes executed from suspicious locations&lt;/p></description></item><item><title>Cisco Modify Configuration</title><link>https://bradleyjkemp.dev/sigmadoc/rules/network/cisco/aaa/cisco_cli_modify_config.yml/</link><pubDate>Thu, 14 Nov 2019 20:55:28 +0100</pubDate><guid>https://bradleyjkemp.dev/sigmadoc/rules/network/cisco/aaa/cisco_cli_modify_config.yml/</guid><description>&lt;p>Modifications to a config that will serve an adversary&amp;rsquo;s impacts or persistence&lt;/p></description></item><item><title>Abusing Windows Telemetry For Persistence</title><link>https://bradleyjkemp.dev/sigmadoc/rules/windows/process_creation/proc_creation_win_abusing_windows_telemetry_for_persistence.yml/</link><pubDate>Mon, 01 Jan 0001 00:00:00 +0000</pubDate><guid>https://bradleyjkemp.dev/sigmadoc/rules/windows/process_creation/proc_creation_win_abusing_windows_telemetry_for_persistence.yml/</guid><description>&lt;p>Windows telemetry makes use of the binary CompatTelRunner.exe to run a variety of commands and perform the actual telemetry collections. This binary was created to be easily extensible, and to that end, it relies on the registry to instruct on which commands to run. The problem is, it will run any arbitrary command without restriction of location or type.&lt;/p></description></item><item><title>Abusing Windows Telemetry For Persistence</title><link>https://bradleyjkemp.dev/sigmadoc/rules/windows/registry_event/registry_event_abusing_windows_telemetry_for_persistence.yml/</link><pubDate>Mon, 01 Jan 0001 00:00:00 +0000</pubDate><guid>https://bradleyjkemp.dev/sigmadoc/rules/windows/registry_event/registry_event_abusing_windows_telemetry_for_persistence.yml/</guid><description>&lt;p>Windows telemetry makes use of the binary CompatTelRunner.exe to run a variety of commands and perform the actual telemetry collections. This binary was created to be easily extensible, and to that end, it relies on the registry to instruct on which commands to run. The problem is, it will run any arbitrary command without restriction of location or type.&lt;/p></description></item><item><title>CrackMapExec Command Execution</title><link>https://bradleyjkemp.dev/sigmadoc/rules/windows/process_creation/proc_creation_win_susp_crackmapexec_execution.yml/</link><pubDate>Mon, 01 Jan 0001 00:00:00 +0000</pubDate><guid>https://bradleyjkemp.dev/sigmadoc/rules/windows/process_creation/proc_creation_win_susp_crackmapexec_execution.yml/</guid><description>&lt;p>Detect various execution methods of the CrackMapExec pentesting framework&lt;/p></description></item><item><title>Exchange Exploitation Activity</title><link>https://bradleyjkemp.dev/sigmadoc/rules/windows/process_creation/proc_creation_win_apt_hafnium.yml/</link><pubDate>Mon, 01 Jan 0001 00:00:00 +0000</pubDate><guid>https://bradleyjkemp.dev/sigmadoc/rules/windows/process_creation/proc_creation_win_apt_hafnium.yml/</guid><description>&lt;p>Detects activity observed by different researchers to be HAFNIUM group activity (or related) on Exchange servers&lt;/p></description></item><item><title>New TaskCache Entry</title><link>https://bradleyjkemp.dev/sigmadoc/rules/windows/registry_event/registry_event_taskcache_entry.yml/</link><pubDate>Mon, 01 Jan 0001 00:00:00 +0000</pubDate><guid>https://bradleyjkemp.dev/sigmadoc/rules/windows/registry_event/registry_event_taskcache_entry.yml/</guid><description>&lt;p>Monitor the creation of a new key under &amp;lsquo;TaskCache&amp;rsquo; when a new scheduled task is registered&lt;/p></description></item></channel></rss>